# roles/to_remove/share-server/templates/check_perm.crontab.j2
*/3 * * * * root {{ glx_share_server_script_path }}/{{ glx_share_server_script_name }}