# roles/powertop/templates/Debian/9/powertop.service.j2
[Unit]
Description=Powertop tunings

[Service]
Type=idle
Environment="TERM=dumb"
ExecStart={{ powertop_bin_path }} --auto-tune
RemainAfterExit=true

[Install]
WantedBy=multi-user.target